# chuhuo_2.71/bd_cli/console调试/etc/init/ttyS0.conf
#ttyS0 - agetty
#
#This service maintains a agetty on ttyS0.
stop on runlevel[S016]
start on runlevel[23]
respawn
exec /sbin/agetty -h -L -w /dev/ttyS0 9600 vt102
