# chuhuo_2.71/sysconf/usr/lib/systemd/system/tog-pegasus.service.d/genssl.conf
[Service]
type=oneshot
ExecStartPre=/usr/share/Pegasus/scripts/generate-certs
