# chuhuo_2.71/sysconf/usr/lib/systemd/system/dirsrv.target
[Unit]
Description=389 Directory Server
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
