# chuhuo_2.71/sysconf/etc/systemd/system/multi-user.target.wants/hypervkvpd.service
[Unit]
Description=Hyper-V KVP daemon
ConditionVirtualization=microsoft

[Service]
Type=simple
ExecStart=/usr/sbin/hypervkvpd

[Install]
WantedBy=multi-user.target
