# chuhuo_2.71/sysconf/etc/systemd/system/multi-user.target.wants/hypervvssd.service
[Unit]
Description=Hyper-V VSS daemon
ConditionVirtualization=microsoft

[Service]
ExecStart=/usr/sbin/hypervvssd

[Install]
WantedBy=multi-user.target
